# package/base-files/files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 74
/svn/ezp/!svn/ver/4553/branches/ray/package/base-files/files/etc/hotplug.d
END
